# benchmark_blueprints/families/codex-surface-workflow-mapping/report/attempt_02_live_probe/workdirs/v4-multi-corpus-objective-run1/workspace/ops/automation/triage.cron
0 9 * * 1-5 python3 scripts/triage.py --window active --emit-md reports/daily_triage.md
